# .local/share/systemd/user/auto-rclone-login.service
[Unit]
Description=Mount OneDrive with rclone on login
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/bin/rclone mount onedrive: %h/onedrive --vfs-cache-mode=writes
ExecStop=/bin/fusermount -u %h/onedrive
Restart=on-failure
TimeoutStopSec=10

[Install]
WantedBy=default.target
